# packages/debian/postinst
#!/bin/bash

# postinst
# post-install script for debian package.
# Creates and configures necessary directories and files and registers daemon.
# Runs post installing/unpacking of the debian package.
# See https://www.debian.org/doc/debian-policy/ch-maintainerscripts.html
#     https://www.debian.org/doc/debian-policy/ap-flowcharts.html

# Any error should exit the script (sort of)
# See https://www.gnu.org/software/bash/manual/html_node/The-Set-Builtin.html#The-Set-Builtin

set -e

adu_conf_dir=/etc/adu
adu_conf_file=du-config.json
adu_diagnostics_conf_file=du-diagnostics-config.json
adu_log_dir=/var/log/adu
adu_data_dir=/var/lib/adu
adu_downloads_dir="$adu_data_dir/downloads"
adu_shell_dir=/usr/lib/adu
adu_shell_file=adu-shell

adu_bin_path=/usr/bin/AducIotAgent
adu_extensions_dir="$adu_data_dir/extensions"
adu_extensions_sources_dir="$adu_extensions_dir/sources"

adu_apt_handler_file=libmicrosoft_apt_1.so
adu_script_handler_file=libmicrosoft_script_1.so
adu_simulator_handler_file=libmicrosoft_simulator_1.so
adu_steps_handler_file=libmicrosoft_steps_1.so
adu_swupdate_handler_file=libmicrosoft_swupdate_1.so

adu_example_component_enumerator_file=libcontoso-component-enumerator.so

adu_delivery_optimization_downloader_file=libdeliveryoptimization-content-downloader.so

adu_eis_conf_file=adu.toml
eis_idservice_dir=/etc/aziot/identityd/config.d

# adu_user is the user that the ADU Agent daemon will run as.
# ADU Agent daemon needs to run as 'adu' user to be able to perform high-privilege tasks via adu-shell.
adu_user=adu

# The adu_group is the group that gives admins and partner users like DO user
# access to ADU resources like sandbox folder.
adu_group=adu

# Use getent and cut to get the info for the adu user and parse out the home directory.
adu_home_dir=$(getent passwd $adu_user | cut -d: -f6)

ms_doclient_lite_service="deliveryoptimization-agent"

# The sample du-config.json
sample_du_config=$(
    cat << END_OF_JSON
{
  "schemaVersion": "1.0",
  "aduShellTrustedUsers": [
    "adu",
    "do"
  ],
  "manufacturer": <Place your device info manufacturer here>,
  "model": <Place your device info model here>,
  "agents": [
    {
      "name": <Place your agent name here>,
      "runas": "adu",
      "connectionSource": {
        "connectionType": "string",
        "connectionData": <Place your Azure IoT device connection string here>
      },
      "manufacturer": <Place your device property manufacturer here>,
      "model": <Place your device property model here>
    }
  ]
}
END_OF_JSON
)

# sample du-diagnostics-config.json.json file
sample_du_diagnostics_config=$(
    cat << END_OF_JSON
{
    "logComponents":[
        {
            "componentName":"adu",
            "logPath":"/var/log/adu/"
        },
        {
            "componentName":"do",
            "logPath":"/var/log/deliveryoptimization-agent/"
        }
    ],
    "maxKilobytesToUploadPerLogPath":50
}
END_OF_JSON
)

setup_dirs_and_files() {
    echo "Setting directories and files..."
    # Note on linux permissions
    # u - The user owner of a file or a directory.
    # g - The group owner of a file or a directory. A group can contain multiple users.
    # o - Any other user that has access to the file or directory.
    # a - All. The combination of user owner, group ownder, and other users.
    # r - Read access.
    #   Read access on a file allows a user to open and read the contents of the file.
    #   Read access on a directory allows a user to list the contents of the directory. e.g with the ls command.
    # w - Write access.
    #   Write access on a file allows a user to modify the contents of a file.
    #   Write access on a directory allows a user to add, remove, rename, or move files in the directory.
    # x - Execute access.
    #   Execute access on a file allows a user to execute the contents of the file as a program.
    #   Execute access on a directory allows a user to enter that directory and possibly gain access to sub-directories. e.g. with the cd command.
    if id -u "$adu_user" > /dev/null 2>&1; then
        mkdir -p "$adu_conf_dir"

        # Ensure the right owners and permissions
        chown "$adu_user:$adu_group" "$adu_conf_dir"
        chmod u=rwx,g=rx,o=rx "$adu_conf_dir"

        # Generate the template configuration file
        echo "Generate the template configuration file..."
        if [ ! -f "$adu_conf_dir/${adu_conf_file}.template" ]; then
            echo "$sample_du_config" > "$adu_conf_dir/${adu_conf_file}.template"
            chown "$adu_user:$adu_group" "$adu_conf_dir/${adu_conf_file}.template"
            chmod u=r,g=r,o=r "$adu_conf_dir/${adu_conf_file}.template"
        fi

        # Create configuration file from template
        if [ ! -f "$adu_conf_dir/$adu_conf_file" ]; then
            cp -a "$adu_conf_dir/${adu_conf_file}.template" "$adu_conf_dir/$adu_conf_file"
            chmod u=rw,g=r,o=r "$adu_conf_dir/$adu_conf_file"
        fi

        echo "Generating the diagnostics configuration file..."
        if [ ! -f "$adu_conf_dir/$adu_diagnostics_conf_file" ]; then
            echo "$sample_du_diagnostics_config" > "$adu_conf_dir/$adu_diagnostics_conf_file"
            chown "$adu_user:$adu_group" "$adu_conf_dir/$adu_diagnostics_conf_file"
            chmod u=r,g=r,o=r "$adu_conf_dir/$adu_diagnostics_conf_file"
        fi

        # Create home dir
        if [ ! -d "$adu_home_dir" ]; then
            mkdir -p "$adu_home_dir"
            chown "$adu_user:$adu_group" "$adu_home_dir"
            chmod u=rwx,g=rx "$adu_home_dir"
        fi

        # Create log dir
        if [ ! -d "$adu_log_dir" ]; then
            mkdir -p "$adu_log_dir"
        fi

        # Ensure the right owners and permissions
        chown "$adu_user:$adu_group" "$adu_log_dir"
        chmod u=rwx,g=rx "$adu_log_dir"

        # Create data dir
        echo "Create data dir..."
        if [ ! -d "$adu_data_dir" ]; then
            mkdir -p "$adu_data_dir"
        fi

        # Ensure the right owners and permissions
        chown "$adu_user:$adu_group" "$adu_data_dir"
        chmod u=rwx,g=rwx,o= "$adu_data_dir"
        ls -ld "$adu_data_dir"

        # Create downloads dir
        if [ ! -d "$adu_downloads_dir" ]; then
            mkdir -p "$adu_downloads_dir"
        fi

        # Ensure the right owners and permissions
        chown "$adu_user:$adu_group" "$adu_downloads_dir"
        chmod u=rwx,g=rwx,o= "$adu_downloads_dir"
        ls -ld "$adu_downloads_dir"

        # Create extensions dir
        if [ ! -d "$adu_extensions_dir" ]; then
            mkdir -p "$adu_extensions_dir"
        fi

        # Ensure the right owners and permissions
        chown "$adu_user:$adu_group" "$adu_extensions_dir"
        chmod u=rwx,g=rwx,o= "$adu_extensions_dir"

        # Create extensions sources dir
        if [ ! -d "$adu_extensions_sources_dir" ]; then
            mkdir -p "$adu_extensions_sources_dir"
        fi

        # Ensure the right owners and permissions
        chown "$adu_user:$adu_group" "$adu_extensions_sources_dir"
        chmod u=rwx,g=rwx,o= "$adu_extensions_sources_dir"

        #Set adu-agent owner and permission
        if [ ! -f "$adu_shell_dir/$adu_shell_file" ]; then
           echo "ERROR! $adu_shell_dir/$adu_shell_file does not exists." >&2
        else
           chown "root:$adu_group" "$adu_shell_dir/$adu_shell_file"
           chmod u=rxs "$adu_shell_dir/$adu_shell_file"
        fi
    else
        echo "ERROR! $adu_user does not exist." >&2
        return 1
    fi
}

register_reference_extensions() {
    echo "Register all reference update content handlers..."

    $adu_bin_path -l 2 --update-type "microsoft/apt:1" -C $adu_extensions_sources_dir/$adu_apt_handler_file
    $adu_bin_path -l 2 --update-type "microsoft/swupdate:1" -C $adu_extensions_sources_dir/$adu_swupdate_handler_file
    $adu_bin_path -l 2 --update-type "microsoft/script:1" -C $adu_extensions_sources_dir/$adu_script_handler_file
    $adu_bin_path -l 2 --update-type "microsoft/steps:1" -C $adu_extensions_sources_dir/$adu_steps_handler_file
    $adu_bin_path -l 2 --update-type "microsoft/update-manifest" -C $adu_extensions_sources_dir/$adu_steps_handler_file
    $adu_bin_path -l 2 --update-type "microsoft/update-manifest:4" -C $adu_extensions_sources_dir/$adu_steps_handler_file

    echo "Register a reference component enumerator..."
    $adu_bin_path -l 2 -E $adu_extensions_sources_dir/$adu_example_component_enumerator_file

    echo "Register content downloader extension..."
    $adu_bin_path -l 2 -D $adu_extensions_sources_dir/$adu_delivery_optimization_downloader_file
}

register_adu_user_and_process_with_eis() {
        if getent group "aziotid" > /dev/null; then
            usermod -aG aziotid $adu_user
        fi

        if getent group "aziotcs" > /dev/null; then
            usermod -aG aziotcs $adu_user
        fi

        if getent group "aziotks" > /dev/null; then
            usermod -aG aziotks $adu_user
        fi

        # Register our process with EIS
        if [ -d "$eis_idservice_dir" ]; then

            if [ ! -f "$eis_idservice_dir/$adu_eis_conf_file" ]; then
                printf "[[principal]]\n name=\"IoTHubDeviceUpdate\"\n idtype=[\"module\"]\n uid= $(id -u $adu_user)\n" > "$eis_idservice_dir/$adu_eis_conf_file"
                chown aziotid:aziotid "$eis_idservice_dir/$adu_eis_conf_file"
                chmod u=rw "$eis_idservice_dir/$adu_eis_conf_file"
            fi
        fi
}

register_and_enable_daemon()
{
    systemctl daemon-reload
    systemctl enable adu-agent
}

case "$1" in
    configure)
        setup_dirs_and_files
        register_adu_user_and_process_with_eis

        register_reference_extensions

        register_and_enable_daemon

        # Note: DO user and group are created by deliveryoptimization-agent Debian package,
        # which is one of the dependencies declared in adu-agent control file.
        # We are assuming that both DO user and group currently exist at this point.
        # Add 'do' user to 'adu' group to allow DO to write to ADU download sandbox.
        echo "Add the 'do' user to the 'adu' group and 'adu' to the 'do' group"
        if getent passwd 'do' > /dev/null; then
            usermod -aG 'adu' 'do'
            # Note: We must add the 'adu' user to the 'do' group so that we can set the connection_string for DO
            usermod -aG 'do' 'adu'

            # Restart deliveryoptimization-agent service to ensure that the new 'do' user's groups take effect.
            echo "Restart $ms_doclient_lite_service"
            systemctl restart "$ms_doclient_lite_service"
        fi

        # Do not restart the current adu-agent because we want the update workflow
        # to complete, and reports the agent's state to ADU Management service.
        echo "============================================================"
        echo
        echo "     ADU Agent service must be restarted." >&2
        echo
        echo "     Please run 'sudo systemctl restart adu-agent'" >&2
        echo
        echo "============================================================"
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

exit 0
